# crontest.cron
15 10 * * *  echo 'date' >>crontest.txt
